# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:redis]
command=redis-server
autostart=true
autorestart=true
stdout_logfile=/tmp/redis.log
stderr_logfile=/tmp/redis.err

[program:api]
command=uvicorn app.main:app --host 127.0.0.1 --port 8000
directory=/app/backend
autostart=true
autorestart=true
environment=REDIS_URL="redis://localhost:6379/0",UPLOAD_DIR="/data/uploads",RESULT_DIR="/data/results"
stdout_logfile=/tmp/api.log
stderr_logfile=/tmp/api.err
depends_on=redis

[program:celery]
command=celery -A app.celery_app:celery_app worker --loglevel=info --concurrency=2
directory=/app/backend
autostart=true
autorestart=true
environment=REDIS_URL="redis://localhost:6379/0",UPLOAD_DIR="/data/uploads",RESULT_DIR="/data/results"
stdout_logfile=/tmp/celery.log
stderr_logfile=/tmp/celery.err
depends_on=redis

[program:nginx]
command=nginx -c /app/nginx.hf.conf
autostart=true
autorestart=true
stdout_logfile=/tmp/nginx.log
stderr_logfile=/tmp/nginx.err
depends_on=api
